# backend/deploy/blog-backend.service
[Unit]
Description=Blog Backend API Service (Express.js)
Documentation=https://github.com/your-username/blog
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Working directory
WorkingDirectory=/home/ubuntu/workspace/blog/backend

# Environment file (contains all env variables)
EnvironmentFile=/home/ubuntu/workspace/blog/backend/.env

# Node.js path (adjust if using nvm or different installation)
ExecStart=/usr/bin/node /home/ubuntu/workspace/blog/backend/src/index.js

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=0

# Process management
KillSignal=SIGINT
KillMode=mixed
TimeoutStopSec=10

# Resource limits (optional)
# LimitNOFILE=65536
# MemoryLimit=512M
# CPUQuota=80%

# Security options (optional but recommended)
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/home/ubuntu/workspace/blog/backend

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blog-backend

# Health check (optional)
# ExecStartPre=/usr/bin/curl -f http://localhost:5080/api/v1/healthz || exit 1

[Install]
WantedBy=multi-user.target
